# dist/src/rpm/resources/prerm
#!/bin/sh -e

echo "Running prerm script..."

if [ $1 -eq 0 ] ; then
    # Package removal, not upgrade
    systemctl --no-reload disable ${packageName}.service > /dev/null 2>&1 || :
    systemctl stop ${packageName}.service > /dev/null 2>&1 || :
fi
